# systemd/lore-auth.service
[Unit]
Description=Lore Auth Broker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lore-auth
Group=lore-auth
ExecStart=/usr/local/bin/lore-auth-server --config /etc/lore-auth/config.yaml
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/lore-auth /etc/lore-auth/keys

[Install]
WantedBy=multi-user.target
